# backup/disk-backup/disk-backup.service
[Unit]
Description=archive root and send to server

[Service]
ExecStart=/usr/local/bin/disk-backup.sh

