# security/tripwire/systemd/exec-tripwire.service
[Unit]
Description=perform an integrity test of Tripware and send it to email

[Service]
ExecStart=/usr/local/bin/exec-tripwire.sh

