# backup/partition-backup/partition-backup.service
[Unit]
Description=archive partitions and send to server

[Service]
ExecStart=/usr/local/bin/partition-backup.sh

