# backup/folder-backup/folder-backup.timer
[Unit]
Description=daily rsync at 12 PM

[Timer]
OnCalendar=Mon..Sun *-*-* 12:00:00
Persistent=True
Unit=folder-backup.service

[Install]
WantedBy=multi-user.target

