# backup/folder-backup/folder-backup.service
[Unit]
Description=rsync local folder with server

[Service]
ExecStart=/usr/local/bin/folder-backup.sh

